# deploy/systemd/atlas-registry-refresh.service
[Unit]
Description=Atlas — refresh atlas.json from on-chain reads
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/beacon
ExecStart=/usr/bin/node /opt/beacon/scripts/build-atlasV2-registry-vps.mjs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=atlas-registry-refresh
