# deploy/systemd/beacon-wallet-risk.service
[Unit]
Description=Beacon signal: wallet-risk
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/beacon/signals/wallet-risk
EnvironmentFile=/opt/beacon/signals/wallet-risk/.env
ExecStartPre=-/bin/fuser -k 4001/tcp
ExecStart=/usr/bin/node /opt/beacon/signals/wallet-risk/dist/index.js
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=beacon-wallet-risk

[Install]
WantedBy=multi-user.target
